# traefik/rootfs/etc/services.d/traefik/run
#!/usr/bin/with-contenv bashio

ENV_VARS=$(gomplate -d options=/data/options.json -i '{{ range (ds "options").env_vars }}{{ . }} {{ end }}')

if [ -z "$ENV_VARS" ]; then
    bashio::log.info "No additional environment variables found"
else
    bashio::log.info "Extracted variables ${ENV_VARS}"
fi

bashio::log.info "Starting Traefik..."
if [ -z "$ENV_VARS" ]; then
    bashio::log.info "Running Traefik without env_vars"
    exec /usr/local/bin/traefik
else
    env $ENV_VARS /usr/local/bin/traefik
fi
